# poky/meta-gui/recipes-hu/hu/files/hustart.service
[Unit]
Description=Persistent NAT bridge and iptables setup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/starthu.sh
Restart=on-failure
RestartSec=3
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
